# collections/ansible_collections/jaredhocutt/homelab/roles/backrest/templates/backrest.service.j2
[Unit]
Description=Backrest
After=network.target

[Service]
Type=simple
User={{ backrest_user }}
Group={{ backrest_group }}
ExecStart={{ backrest_bin_path }}
{% for key, value in (common_container_env_defaults | combine(backrest_container_env_defaults, recursive=True) | combine(backrest_container_env, recursive=True)).items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}

[Install]
WantedBy=multi-user.target
